# niri/sessions/dms.service.d/combo.conf
# Gate DMS so it auto-starts in every niri session EXCEPT the explicit
# "niri + Noctalia" combo. The session wrapper exports NIRI_SHELL and
# niri-session imports it into the systemd --user manager before niri.service
# (and thus this WantedBy=niri.service unit) starts.
#
# Negated match: start DMS unless NIRI_SHELL == noctalia. This keeps the stock
# "Niri" session and the "niri + DMS" session working as before; only the
# Noctalia combo skips DMS. Requires systemd >= 246 (ConditionEnvironment).
[Unit]
ConditionEnvironment=!NIRI_SHELL=noctalia
